# src/etc/supervisor.d/start_crond.conf
[program:crond]
command=crond -f
stdout_events_enabled = true
stderr_events_enabled = true
autostart=true