# src/etc/supervisor.d/start_php-fpm.conf
[program:php-fpm]
command=php-fpm8 -y /etc/php8/php-fpm.conf --nodaemonize
stdout_events_enabled = true
stderr_events_enabled = true
autostart=true